# .timers/docker_cleanup.service
[Unit]
Description=Remove dangling docker volumes

[Service]
User=san
ExecStart=/home/san/.scripts/docker_cleanup.sh
