# config/systemd/user/backup.service
[Unit]
Description=Save important files to off-world server

[Service]
User=san
ExecStart=/home/san/.scripts/backup.sh
