# config/systemd/user/cache_cleanup.service
[Unit]
Description=Remove versions older than latest in pacman cache

[Service]
User=root
ExecStart=/home/san/.scripts/clear_pacman_cache.sh
